# installer/ansible/roles/upgrade/templates/odio-upgrade.service.j2
[Unit]
Description=Apply pending odios upgrade
Documentation=man:odio-upgrade(1)
# No safe mid-run cancel: refuse manual stop rather than fake success.
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=/usr/bin/sudo /usr/local/bin/odio-upgrade apply --progress{% if upgrade_pinned_version %} --version {{ upgrade_pinned_version }}{% endif %}
