# installer/ansible/roles/upgrade/templates/odio-check-upgrade.timer.j2
[Unit]
Description=Daily odio upgrade check

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
