# installer/ansible/roles/upgrade/templates/odio-check-upgrade.service.j2
[Unit]
Description=Check for odio upgrades

[Service]
Type=oneshot
ExecStart=/usr/local/bin/odio-upgrade check
