# installer/ansible/roles/upgrade/templates/odio-upgrade.service.j2
[Unit]
Description=Apply pending odios upgrade
Documentation=man:odio-upgrade(1)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/odio-upgrade
